# installer/templates/deft.service
[Unit]
Description=Deft Node Agent
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/deft serve
Restart=always
RestartSec=5
StartLimitIntervalSec=0
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deft

[Install]
WantedBy=multi-user.target
